# packages/my-great-rest-service/redhat/SOURCES/my-great-rest-service.service
[Unit]
Description=My Great Rest Service

[Service]
PermissionsStartOnly=true
User=xroad-catalog
Group=xroad-catalog
Restart=on-failure
ExecStart=/usr/share/xroad/bin/my-great-rest-service

[Install]
WantedBy=multi-user.target




